# contrib/init/coordinated.conf
description "Bitcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env COORDINATED_BIN="/usr/bin/coordinated"
env COORDINATED_USER="bitcoin"
env COORDINATED_GROUP="bitcoin"
env COORDINATED_PIDDIR="/var/run/coordinated"
# upstart can't handle variables constructed with other variables
env COORDINATED_PIDFILE="/var/run/coordinated/coordinated.pid"
env COORDINATED_CONFIGFILE="/etc/bitcoin/bitcoin.conf"
env COORDINATED_DATADIR="/var/lib/coordinated"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # coordinated will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$COORDINATED_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run coordinated."
        echo "The setting must appear in $COORDINATED_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$COORDINATED_PIDDIR"
    chmod 0755 "$COORDINATED_PIDDIR"
    chown $COORDINATED_USER:$COORDINATED_GROUP "$COORDINATED_PIDDIR"
    chown $COORDINATED_USER:$COORDINATED_GROUP "$COORDINATED_CONFIGFILE"
    chmod 0660 "$COORDINATED_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$COORDINATED_PIDFILE" \
    --chuid $COORDINATED_USER:$COORDINATED_GROUP \
    --exec "$COORDINATED_BIN" \
    -- \
    -pid="$COORDINATED_PIDFILE" \
    -conf="$COORDINATED_CONFIGFILE" \
    -datadir="$COORDINATED_DATADIR" \
    -disablewallet \
    -daemon

